# docs/operations/launchd/com.missioncontrol.snapshot-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
<key>Label</key><string>com.missioncontrol.snapshot-publish</string>
<key>ProgramArguments</key><array><string>/bin/bash</string><string>/Users/REPLACE_ME/.openclaw/workspace/repos/ai-clawbot-mission-control/scripts/publish_snapshot.sh</string></array>
<key>WorkingDirectory</key><string>/Users/REPLACE_ME/.openclaw/workspace/repos/ai-clawbot-mission-control</string>
<key>StartInterval</key><integer>300</integer>
<key>RunAtLoad</key><true/>
<key>EnvironmentVariables</key><dict><key>SAFE_MODE</key><string>true</string><key>SNAPSHOT_REPO</key><string>developwithJB/ai-clawbot-mission-control</string><key>SNAPSHOT_BRANCH</key><string>main</string><key>SNAPSHOT_PAT</key><string>REPLACE_WITH_PAT</string></dict>
<key>StandardOutPath</key><string>/tmp/mission-control-snapshot.out.log</string>
<key>StandardErrorPath</key><string>/tmp/mission-control-snapshot.err.log</string>
</dict></plist>
